# package/vpnpool/files/etc/init.d/vpnpool
#!/bin/sh /etc/rc.common
# vpnpool procd service (stage 1 skeleton).
# Supervises the control-plane daemon. Does NOT touch podkop/zapret/routing yet.

USE_PROCD=1
START=99
STOP=10

PROG=/usr/libexec/vpnpool/vpnpoold
TGBOT=/usr/libexec/vpnpool/tgbot.sh

# Telegram control bot as its OWN procd instance, so procd is its single supervisor
# (reliable stop/restart, no orphaned long-poller, no pidfile races — see the note in
# vpnpoold). Conditional on telegram_control=1 + token + chat. Kept independent of the
# main `enabled` toggle so the bot can still receive a remote /on when routing is off.
start_tgbot_instance() {
	local ctl tok chat
	config_get_bool ctl main telegram_control 0
	config_get tok main telegram_token
	config_get chat main telegram_chat
	[ "$ctl" = "1" ] && [ -n "$tok" ] && [ -n "$chat" ] || return 0
	procd_open_instance tgbot
	procd_set_param command "$TGBOT"
	procd_set_param respawn 3600 5 0
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance
}

start_service() {
	config_load vpnpool

	local enabled
	config_get_bool enabled main enabled 0
	if [ "$enabled" = "1" ]; then
		procd_open_instance vpnpoold
		procd_set_param command "$PROG"
		procd_set_param respawn 3600 5 0
		procd_set_param stdout 1
		procd_set_param stderr 1
		procd_set_param pidfile /var/run/vpnpool.pid
		procd_close_instance
	else
		logger -t vpnpool "disabled in /etc/config/vpnpool (option enabled '0'); daemon not starting"
	fi

	start_tgbot_instance
}

# On a config-change reload procd re-runs start_service and diffs instances: the daemon
# instance is left running if its params are unchanged, while the tgbot instance is
# started/stopped to match telegram_control. So toggling the bot in the UI never bounces
# the tunnel. rpcd's set_option calls `/etc/init.d/vpnpool reload` for telegram_control.
reload_service() {
	start_service
}

stop_service() {
	# Routing/sing-box teardown is done by vpnpoold's own TERM handler; procd stops both
	# instances (daemon + tgbot) on stop/restart.
	return 0
}

service_triggers() {
	procd_add_reload_trigger "vpnpool"
}
